# xdg/etc/profile.d/xdg-gnupg.sh
# GnuPG XDG Specification

export GNUPGHOME="${XDG_DATA_HOME}"/gnupg
