# xdg/etc/profile.d/xdg-zsh.sh
# zsh XDG Specification

export ZDOTDIR="${XDG_CONFIG_HOME}"/zsh
