# xdg/etc/profile.d/xdg-bash.sh
# BASH XDG Specification

. "${XDG_CONFIG_HOME}"/bash/bash_profile
. "${XDG_CONFIG_HOME}"/bash/bashrc

export HISTFILE="${XDG_DATA_HOME}"/bash/history
