# qemu-kvm/systemd/network-bridge.service
[Unit]
Description=Setup a Network Bridge Using iproute2

[Service]
ExecStart=/usr/local/sbin/setup_br0.sh

[Install]
WantedBy=multi-user.target 
